# etc-udev-rules_d/92-keyboard-hotplug.rules
SUBSYSTEM=="input", SUBSYSTEMS=="usb", ACTION=="change", ENV{XAUTHORITY}="/home/motoko/.Xauthority", RUN+="/bin/bash /home/motoko/.config/mykeymaps/xmodmap.sh"

